# deployment/systemd/portfolio-nextjs.service
[Unit]
Description=The Tirtayasa Portfolio Next.js frontend
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/thetirtayasa/frontend
Environment=NODE_ENV=production
Environment=HOSTNAME=127.0.0.1
Environment=PORT=3030
ExecStart=/usr/bin/node .next/standalone/server.js
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
